# scripts/internal/DO_NOT_TOUCH_HUMAN/rpi_ai.service
[Unit]
Description=RPI AI Agent Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=davca
WorkingDirectory=/home/davca/rpi_ai/rpi_ai
ExecStart=/home/davca/rpi_ai/rpi_ai/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
